# initscript.conf
#/etc/init/salink.conf

description "uWSGI server instance configured to serve collabdesk"

start on runlevel [2345]
stop on runlevel [!2345]

setuid waqas
setgid www-data

env PATH=/home/waqas/salink/env/bin
chdir /home/waqas/salink
exec uwsgi --ini salink.ini
