# initscript.conf
#/etc/init/pf.conf

description "uWSGI server instance configured to serve collabdesk"

start on runlevel [2345]
stop on runlevel [!2345]

setuid waqas
setgid www-data

env PATH=/home/waqas/pf/env/bin
chdir /home/waqas/pf
exec uwsgi --ini pf.ini
